# docker/dev-monolith/s6-overlay/s6-rc.d/redis-log/run
#!/bin/sh
exec logutil-service /var/log/redis