# docker/dev-monolith/s6-overlay/s6-rc.d/cockroachdb-health/up
/etc/s6-overlay/scripts/cockroachdb-healthcheck-loop.sh