# docker/dev-monolith/s6-overlay/s6-rc.d/cockroachdb/run
#!/bin/sh
exec 2>&1
exec su cockroachdb -c 'cockroach start-single-node --insecure --store=/data/cockroachdb --http-port 9200'