# docker/dev-monolith/s6-overlay/s6-rc.d/clickhouse/run
#!/bin/sh
exec 2>&1
exec su clickhouse -c 'clickhouse-server --config-file=/etc/clickhouse-server/config.xml'