# docker/dev-monolith/s6-overlay/s6-rc.d/rivet-server/run
#!/bin/sh
exec 2>&1
exec su rivet-server -c 'rivet-server start'