# docker/dev-monolith/s6-overlay/s6-rc.d/rivet-server-health/up
/etc/s6-overlay/scripts/rivet-server-healthcheck-loop.sh