# docker/dev-monolith/s6-overlay/s6-rc.d/clickhouse-health/up
/etc/s6-overlay/scripts/clickhouse-healthcheck-loop.sh