# docker/dev-monolith/s6-overlay/s6-rc.d/rivet-client/run
#!/bin/sh
exec 2>&1
exec rivet-client -c /etc/rivet-client/config.yaml