# docker/dev-monolith/s6-overlay/s6-rc.d/redis/run
#!/bin/sh
exec 2>&1
exec su redis -c 'redis-server --dir /data/redis --requirepass password --save 60 1 --appendonly yes'