# docker/dev-monolith/s6-overlay/s6-rc.d/redis-health/up
/etc/s6-overlay/scripts/redis-healthcheck-loop.sh