# docker/dev-monolith/s6-overlay/s6-rc.d/nats/run
#!/bin/sh
exec 2>&1
exec su nats -c 'nats-server'