# docker/dev-monolith/s6-overlay/s6-rc.d/seaweedfs-health/up
/etc/s6-overlay/scripts/seaweedfs-healthcheck-loop.sh