# docker/dev-monolith/s6-overlay/s6-rc.d/nats-health/up
/etc/s6-overlay/scripts/nats-healthcheck-loop.sh