# docker/dev-monolith/s6-overlay/s6-rc.d/rivet-client-log/run
#!/bin/sh
exec logutil-service /var/log/rivet-client